# scripts/xero-sync.service
[Unit]
Description=Xero Sync Service for %i
After=network.target

[Service]
Type=simple
User=%i
Group=%i
EnvironmentFile=/usr/local/etc/xero-sync.env
WorkingDirectory=${WORKDIR}
ExecStart=${VENV_PYTHON} ${WORKDIR}/manage.py start_xero_sync
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
